# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/run/supervisord.pid
user=root

[program:hub-api]
command=python3 /opt/hub/docker/hub-api.py
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DATA_ROOT="/data",HUB_API_PORT="7870"

[program:hub-bootstrap]
command=/bin/bash -c 'sleep 15 && /opt/hub/docker/switch-app.sh sillytavern'
autostart=true
autorestart=false
startretries=3
priority=8
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DATA_ROOT="/data",ACTIVE_APP="%(ENV_ACTIVE_APP)s",ST_PORT="8000",LUMIVERSE_PORT="7861",MARINARA_PORT="7862"

[program:nginx]
command=/usr/sbin/nginx -c /opt/hub/docker/nginx.conf -g "daemon off;"
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:sillytavern]
command=/opt/hub/docker/run-sillytavern.sh
directory=/apps/sillytavern
user=node
autostart=false
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DATA_ROOT="/data",ST_PORT="8000"

[program:lumiverse]
command=/opt/hub/docker/run-lumiverse.sh
directory=/apps/lumiverse
user=root
autostart=false
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DATA_ROOT="/data",LUMIVERSE_PORT="7861"

[program:marinara]
command=/opt/hub/docker/run-marinara.sh
directory=/apps/marinara
user=node
autostart=false
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DATA_ROOT="/data",MARINARA_PORT="7862"

[program:sync-shared]
command=/bin/bash -c 'while true; do /opt/hub/scripts/sync-shared-data.sh || true; sleep 300; done'
autostart=true
autorestart=true
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0